# contrib/idioticio.service
[Service]
Type=simple
ExecStart=/usr/bin/idioticio /etc/idiotic.io/conf.yaml
